# templates/uwsgi_systemd.service.erb
#
# This file is managed by the puppet class 'uwsgi'
#
# uWSGI Emperor process systemd service
#
[Unit]
Description=uWSGI Emperor
After=syslog.target

[Service]
ExecStart=<%= @binary_directory %>/uwsgi --ini <%= @config_file %>
# Requires systemd version 211 or newer
RuntimeDirectory=uwsgi
Restart=always
KillSignal=SIGQUIT
Type=notify
StandardError=syslog
NotifyAccess=all

[Install]
WantedBy=multi-user.target
