# templates/uwsgi_upstart.conf.erb
#
# This file is managed by puppet class 'uwsgi'
#
# uWSGI Emperor process Upstart script
#
description "uWSGI Emperor"
start on runlevel [2345]
stop on runlevel [!2345]

respawn

pre-start script
    uwsgipid="<%= @pidfile %>"
    uwsgisocket="<%= @socket %>"
    uwsgipiddir="${uwsgipid%/*}"
    uwsgisocketdir="${uwsgisocket%/*}"
    mkdir -p "$uwsgipiddir"
    mkdir -p "$uwsgisocketdir"
    mkdir -p /var/log/uwsgi
end script
exec uwsgi --die-on-term --ini <%= @config_file %>
