# templates/uwsgi_systemd.service.erb
#
# This file is managed by the puppet class 'uwsgi'
#
# uWSGI Emperor process systemd service
#
[Unit]
Description=uWSGI Emperor
After=syslog.target

[Service]
ExecStart=<%= @binary_directory %>/uwsgi <%= @die_on_term %> --ini <%= @configfile %>
# Requires systemd version 211 or newer
RuntimeDirectory=uwsgi
Restart=always
KillSignal=<%= @kill_signal %>
Type=notify
StandardError=syslog
NotifyAccess=all

[Install]
WantedBy=multi-user.target
